# jobs/mongodb_node/templates/services.conf
pre-start script
  chown vcap:vcap -R /store
end script

#start on runlevel [2345]
start on filesystem
stop on runlevel S
